# systemd/offlineimap/offlineimap.service
[Unit]
Description=offlineimap sync service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -o -u basic
ExecStop=/usr/bin/pkill offlineimap
TimeoutSec=100
TimeoutStopSec=120
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
