# systemd/offlineimap/offlineimap.target
[Unit]
Description=offlineimap target
StopWhenUnneeded=yes
