# systemd/offlineimap/offlineimap.timer
[Unit]
Description=2 minute offlineimap timer

[Timer]
OnBootSec=2min
OnUnitInactiveSec=15m
Persistent=yes

[Install]
WantedBy=default.target
