# systemd/offlineimap/offlineimap.timer
[Unit]
Description=2 minute offlineimap timer

[Timer]
OnBootSec=2min
OnCalendar=*-*-* *:0/2
Persistent=yes

[Install]
WantedBy=timers.target
