# systemd/offlineimap/offlineimap.service
[Unit]
Description=offlineimap sync service

[Service]
Type=simple
ExecStart=/usr/bin/offlineimap -o -u basic
TimeoutStopSec=60
; Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
