# systemd/bup/bup.service
[Unit]
Description = Bups backup manager service

[Service]
Type = simple

