# systemd/bup/bup.timer
[Unit]
Description = backup timer

[Timer]
OnCalendar = *-*-*
Persistent = true

[Install]
WantedBy = timers.target

