# sudo/offlineimap.timer
[Unit]
Description=1 Minute timer

[Timer]
OnCalendar=*-*-* *:*:34
