# sudo/offlineimap.service
[Unit]
Description=offlineimap sync service

[Service]
Type=oneshot
ExecStart=/usr/vin/offlineimap
