# sudo/offlineimap.timer
[Unit]
Description=2 Minute timer

[Timer]
OnBootSec=1min
OnCalendar=*-*-* *:*:34
Unit=offlineimap-timer.target

[Install]
WantedBy=
